# extras/crontab
# point to path with the script...
*/30 * * * * /home/jmarkow/restore-downed-nodes.sh
